# systemd/overwitch.service
[Unit]
Description=Overwitch service
After=pipewire.service
Requires=pipewire.service

[Service]
GuessMainPID=true
ExecStart=overwitch-service
Restart=on-failure

[Install]
WantedBy=default.target
